# openldap-sasl/files/ldap.bak
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2011 OpenWrt.org

START=51

PIDFILE=/var/run/slapd.pid
CONFFILE=/var/etc/slapd.conf
CLNTCONF=/var/etc/ldap.conf
DATADIR=/overlay/etc/openldap/data
SUFFIX=.ldif
NSSWITCH=/var/etc/nsswitch.conf
BIN=/usr/sbin/slapd

add_line() {
	echo >>$CONFFILE -e $1
}

add_repl() {
	config_get uri $1 uri ""
	[ ! -z "$uri" ] && {
		add_line ""
		add_line "syncRepl rid=$rid"
		add_line "\tprovider=$uri"
		add_line "\tbinddn=\"cn=Manager,ou=System Users,dc=ta,dc=net\""
		add_line "\tbindmethod=simple"
		add_line "\tcredentials=12345"
		add_line "\tsearchbase=$basedn"
		add_line "\ttype=refreshAndPersist"
	}
	mirrormode=true
}

start_server() {
        config_get enabled  $1 enabled   "0"
	config_get serverid $1 serverid "1"
	config_get basedn   $1 basedn   "dc=local,dc=net"

	[ "$enabled" -eq "1" ] && {
		mkdir -p /var/etc
		mkdir -p /var/run
		sed -e "s#|BASEDN|#$basedn#g" \
		    -e "s#|DATADIR|#$DATADIR#g" \
		    -e "s#|SERVERID|#$serverid#g" \
		    /etc/openldap/slapd.conf.template > $CONFFILE
		chmod 600 $CONFFILE

		rid=1
		config_foreach  add_repl replication
		[ ! -z "$mirrormode" ] && add_line "\nmirrormode true"

		$BIN -h "ldapi:/// ldap://" -f $CONFFILE

		[ -e /var/run/ldapi ] && {
			echo -e "URI\tldapi:///" >$CLNTCONF
			chmod 644 $CLNTCONF
			sed -i -e "s/passwd:files$/passwd:files ldap/g" -e "s/group:files$/group:files ldap/g" $NSSWITCH
		}
	}
}

start() {
	config_load ldap
	config_foreach  start_server server
}

stop() {
	[ -f $PIDFILE ] && {
		sed -i -e "s/passwd:files ldap$/passwd:files/g" -e "s/group:files ldap$/group:files/g" $NSSWITCH
		kill $(cat $PIDFILE)
	}
}
