# netatalk3/files/afpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2012 OpenWrt.org

START=70
STOP=80

start()
{
	mkdir -p /var/netatalk/msg
	mkdir -p /var/netatalk/CNID
	service_start /usr/sbin/cnid_metad
	service_start /usr/sbin/afpd
}

stop()
{
	service_stop /usr/sbin/afpd
	service_stop /usr/sbin/cnid_metad
}

reload()
{
	service_reload /usr/sbin/afpd
}
