# xPL4Linux/files/xplhub.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2012 OpenWrt.org
START=80
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/xplhub -nodaemon
	procd_set_param respawn
	procd_close_instance 
}
